# src/a2_unitree_agent/systemd/a2-unitree-agent.service
[Unit]
Description=A2 Unitree SDK isolation agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=A2_WORKSPACE=/home/unitree/ws/device-navigation
Environment=A2_SDK_INTERFACE=eth0
Environment=A2_UNITREE_AGENT_SOCKET=/run/a2/unitree_agent.sock
Environment=LD_LIBRARY_PATH=/opt/unitree_robotics/lib/x86_64:/opt/unitree_robotics/lib
ExecStart=/home/unitree/ws/device-navigation/install/a2_unitree_agent/share/a2_unitree_agent/scripts/start_unitree_agent.sh
Restart=always
RestartSec=1
RuntimeDirectory=a2
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
